# .service
[Unit]
Description=Postgres Database Service
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Restart=always
ExecStart=docker run --name postgres-db -e POSTGRES_PASSWORD=$POSTGRES_PASSWORD -d postgres:15.2

[Install]
WantedBy=default.target
