# .service
[Unit]
Description=Postgres Database Service
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Restart=always
ExecStart=docker start postgres-db
User=pi

[Install]
WantedBy=multi-user.target
